# KF5/kf5-kglobalaccel/kglobalaccel.service
[Unit]
Description=KGlobalaccel

[Service]
Type=dbus
BusName=org.kde.kglobalaccel
ExecStart=/usr/bin/kglobalaccel5
